# bh.socket
[Socket]
ListenFIFO=%t/bh.stdin
Service=bh.service
